# tvserver_server.conf
# Supervisord configuration to run tv server client at boot
[program:tvserver-server]
command=node /home/pi/tvwall/server.js
autostart=true
autorestart=true
startsecs=5