# tvserver_client.conf
# Supervisord configuration to run tv server client at boot
[program:tvserver-client]
command=node /home/pi/tvwall/client.js
autostart=true
autorestart=true
startsecs=5
